# root/etc/systemd/system/backup.automount
[Unit]
Description=Automount NFS share for backup

[Automount]
Where=/backup
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
