# home/.config/systemd/user/backup.service
[Unit]
Description=Back up

[Service]
Type=oneshot
ExecStart=/home/user/backup.sh
